# templates/nginx.socket.j2
[Socket]
ListenStream={{ openslides_nginx_port }}
BindIPv6Only=both

[Install]
WantedBy=sockets.target
