# templates/openslides.socket.j2
[Socket]
ListenStream={{ openslides_port }}

[Install]
WantedBy=sockets.target
